# provisioning/templates/rmq-mkeapi.conf
[program:mkeapi]
command=/usr/bin/php /opt/mq-resources/mkeapi-send-form.php -l
numprocs=1
process_name=send-form-%(process_num)02d
autostart=true
autorestart=true
killasgroup=true
