# provisioning/templates/rmq-recipes.conf
[program:recipes]
command=/usr/bin/php /opt/mq-resources/toh-update-recipe.php -c vagrant
numprocs=1
process_name=toh-%(process_num)02d
autostart=true
autorestart=true
killasgroup=true
